# com.prodigy.wins-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.prodigy.wins-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kenny/Documents/Claude/prodigyathletics-site/scripts/wins_watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kenny/Library/Logs/prodigy_wins_watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kenny/Library/Logs/prodigy_wins_watcher.log</string>
</dict>
</plist>
